# systemd/.config/systemd/user/check-battery-user.service
[Unit]
Description=check-battery

[Service]
ExecStart=/home/meir/.config/scripts/check_battery.sh
